# change_mac.service
[Unit]
Description=Change MAC Address Service

[Service]
ExecStart=CHANGE_TO_YOUR_SHELL_PATH/change_mac.sh enp1s0
