# change_mac.service
[Unit]
Description=Change MAC Address Service

[Service]
ExecStart=/home/devops/Documents/change_mac.sh enp1s0
